# aero-optical-flow.service.in
[Unit]
Description=Aero optical flow

[Service]
Type=simple
ExecStart=/usr/bin/aero-optical-flow
Restart=on-failure

[Install]
WantedBy=multi-user.target
